# 4.Rootfs/var/lib/dpkg/info/bsdmainutils.list
/.
/etc
/etc/calendar
/etc/calendar/default
/etc/cron.daily
/etc/cron.daily/bsdmainutils
/etc/default
/etc/default/bsdmainutils
/usr
/usr/bin
/usr/bin/bsd-from
/usr/bin/bsd-write
/usr/bin/calendar
/usr/bin/col
/usr/bin/colcrt
/usr/bin/colrm
/usr/bin/column
/usr/bin/hexdump
/usr/bin/look
/usr/bin/lorder
/usr/bin/ncal
/usr/bin/printerbanner
/usr/bin/ul
/usr/share
/usr/share/calendar
/usr/share/calendar/calendar.all
/usr/share/calendar/calendar.argentina
/usr/share/calendar/calendar.australia
/usr/share/calendar/calendar.belgium
/usr/share/calendar/calendar.birthday
/usr/share/calendar/calendar.christian
/usr/share/calendar/calendar.computer
/usr/share/calendar/calendar.croatian
/usr/share/calendar/calendar.debian
/usr/share/calendar/calendar.discordian
/usr/share/calendar/calendar.dutch
/usr/share/calendar/calendar.eu
/usr/share/calendar/calendar.freebsd
/usr/share/calendar/calendar.french
/usr/share/calendar/calendar.german
/usr/share/calendar/calendar.history
/usr/share/calendar/calendar.holiday
/usr/share/calendar/calendar.hungarian
/usr/share/calendar/calendar.judaic.2018
/usr/share/calendar/calendar.judaic.2019
/usr/share/calendar/calendar.judaic.2020
/usr/share/calendar/calendar.judaic.2021
/usr/share/calendar/calendar.judaic.2022
/usr/share/calendar/calendar.kazakhstan
/usr/share/calendar/calendar.lotr
/usr/share/calendar/calendar.music
/usr/share/calendar/calendar.newzealand
/usr/share/calendar/calendar.pagan
/usr/share/calendar/calendar.russian
/usr/share/calendar/calendar.southafrica
/usr/share/calendar/calendar.thai
/usr/share/calendar/calendar.ubuntu
/usr/share/calendar/calendar.ukrainian
/usr/share/calendar/calendar.unitedkingdom
/usr/share/calendar/calendar.usholiday
/usr/share/calendar/calendar.world
/usr/share/calendar/de_AT
/usr/share/calendar/de_AT/calendar.feiertag
/usr/share/calendar/de_DE
/usr/share/calendar/de_DE/calendar.all
/usr/share/calendar/de_DE/calendar.feiertag
/usr/share/calendar/de_DE/calendar.geschichte
/usr/share/calendar/de_DE/calendar.kirche
/usr/share/calendar/de_DE/calendar.literatur
/usr/share/calendar/de_DE/calendar.musik
/usr/share/calendar/de_DE/calendar.wissenschaft
/usr/share/calendar/fr_FR
/usr/share/calendar/fr_FR/calendar.all
/usr/share/calendar/fr_FR/calendar.fetes
/usr/share/calendar/fr_FR/calendar.french
/usr/share/calendar/fr_FR/calendar.jferies
/usr/share/calendar/fr_FR/calendar.proverbes
/usr/share/calendar/hr_HR
/usr/share/calendar/hr_HR/calendar.all
/usr/share/calendar/hr_HR/calendar.praznici
/usr/share/calendar/hu_HU
/usr/share/calendar/hu_HU/calendar.all
/usr/share/calendar/hu_HU/calendar.nevnapok
/usr/share/calendar/hu_HU/calendar.unnepek
/usr/share/calendar/pt_BR
/usr/share/calendar/pt_BR/calendar.all
/usr/share/calendar/pt_BR/calendar.commemorative
/usr/share/calendar/pt_BR/calendar.holidays
/usr/share/calendar/pt_BR/calendar.mcommemorative
/usr/share/calendar/ru_RU
/usr/share/calendar/ru_RU/calendar.all
/usr/share/calendar/ru_RU/calendar.common
/usr/share/calendar/ru_RU/calendar.holiday
/usr/share/calendar/ru_RU/calendar.military
/usr/share/calendar/ru_RU/calendar.orthodox
/usr/share/calendar/ru_RU/calendar.pagan
/usr/share/calendar/uk_UA
/usr/share/calendar/uk_UA/calendar.all
/usr/share/calendar/uk_UA/calendar.holiday
/usr/share/calendar/uk_UA/calendar.misc
/usr/share/calendar/uk_UA/calendar.orthodox
/usr/share/doc
/usr/share/doc/bsdmainutils
/usr/share/doc/bsdmainutils/README
/usr/share/doc/bsdmainutils/calendarJudaic.py.gz
/usr/share/doc/bsdmainutils/changelog.armel.gz
/usr/share/doc/bsdmainutils/changelog.gz
/usr/share/doc/bsdmainutils/copyright
/usr/share/doc/bsdmainutils/source.data.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/bsdmainutils
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/bsd-from.1.gz
/usr/share/man/man1/bsd-write.1.gz
/usr/share/man/man1/calendar.1.gz
/usr/share/man/man1/col.1.gz
/usr/share/man/man1/colcrt.1.gz
/usr/share/man/man1/colrm.1.gz
/usr/share/man/man1/column.1.gz
/usr/share/man/man1/hexdump.1.gz
/usr/share/man/man1/look.1.gz
/usr/share/man/man1/lorder.1.gz
/usr/share/man/man1/ncal.1.gz
/usr/share/man/man1/printerbanner.1.gz
/usr/share/man/man1/ul.1.gz
/usr/bin/cal
/usr/bin/hd
/usr/share/calendar/calendar.judaic
/usr/share/man/man1/cal.1.gz
/usr/share/man/man1/hd.1.gz
