# 4.Rootfs/var/lib/dpkg/info/openssh-server.conffiles
/etc/default/ssh
/etc/init.d/ssh
/etc/pam.d/sshd
/etc/ssh/moduli
/etc/ufw/applications.d/openssh-server
