# 4.Rootfs/var/lib/dpkg/info/rsyslog.list
/.
/etc
/etc/default
/etc/default/rsyslog
/etc/init.d
/etc/init.d/rsyslog
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/rsyslog
/etc/logrotate.d
/etc/logrotate.d/rsyslog
/etc/rsyslog.conf
/etc/rsyslog.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rsyslog.service
/usr
/usr/lib
/usr/lib/arm-linux-gnueabi
/usr/lib/arm-linux-gnueabi/rsyslog
/usr/lib/arm-linux-gnueabi/rsyslog/fmhash.so
/usr/lib/arm-linux-gnueabi/rsyslog/imfile.so
/usr/lib/arm-linux-gnueabi/rsyslog/imjournal.so
/usr/lib/arm-linux-gnueabi/rsyslog/imklog.so
/usr/lib/arm-linux-gnueabi/rsyslog/imkmsg.so
/usr/lib/arm-linux-gnueabi/rsyslog/immark.so
/usr/lib/arm-linux-gnueabi/rsyslog/impstats.so
/usr/lib/arm-linux-gnueabi/rsyslog/imptcp.so
/usr/lib/arm-linux-gnueabi/rsyslog/imtcp.so
/usr/lib/arm-linux-gnueabi/rsyslog/imudp.so
/usr/lib/arm-linux-gnueabi/rsyslog/imuxsock.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmnet.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmnetstrms.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmnsd_ptcp.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmregexp.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmtcpclt.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmtcpsrv.so
/usr/lib/arm-linux-gnueabi/rsyslog/lmzlibw.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmanon.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmexternal.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmfields.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmjsonparse.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmnormalize.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmpstrucdata.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmrm1stspace.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmsequence.so
/usr/lib/arm-linux-gnueabi/rsyslog/mmutf8fix.so
/usr/lib/arm-linux-gnueabi/rsyslog/omjournal.so
/usr/lib/arm-linux-gnueabi/rsyslog/ommail.so
/usr/lib/arm-linux-gnueabi/rsyslog/omprog.so
/usr/lib/arm-linux-gnueabi/rsyslog/omuxsock.so
/usr/lib/arm-linux-gnueabi/rsyslog/pmaixforwardedfrom.so
/usr/lib/arm-linux-gnueabi/rsyslog/pmcisconames.so
/usr/lib/arm-linux-gnueabi/rsyslog/pmlastmsg.so
/usr/lib/arm-linux-gnueabi/rsyslog/pmsnare.so
/usr/lib/rsyslog
/usr/lib/rsyslog/rsyslog-rotate
/usr/sbin
/usr/sbin/rsyslogd
/usr/share
/usr/share/doc
/usr/share/doc/rsyslog
/usr/share/doc/rsyslog/AUTHORS
/usr/share/doc/rsyslog/NEWS.Debian.gz
/usr/share/doc/rsyslog/README.Debian
/usr/share/doc/rsyslog/changelog.Debian.gz
/usr/share/doc/rsyslog/changelog.gz
/usr/share/doc/rsyslog/copyright
/usr/share/doc/rsyslog/examples
/usr/share/doc/rsyslog/examples/rsyslog.d
/usr/share/doc/rsyslog/examples/rsyslog.d/console.conf
/usr/share/doc/rsyslog/examples/rsyslog.d/xconsole.conf
/usr/share/doc/rsyslog/examples/tmpfiles.d
/usr/share/doc/rsyslog/examples/tmpfiles.d/xconsole.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rsyslog
/usr/share/man
/usr/share/man/man1
/usr/share/man/man5
/usr/share/man/man5/rsyslog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/rsyslogd.8.gz
/var
/var/spool
/var/spool/rsyslog
