# rootfs/var/lib/dpkg/info/bsdmainutils.conffiles
/etc/calendar/default
/etc/cron.daily/bsdmainutils
/etc/default/bsdmainutils
