# rootfs/var/lib/dpkg/info/cron.list
/.
/etc
/etc/cron.d
/etc/cron.d/.placeholder
/etc/cron.daily
/etc/cron.daily/.placeholder
/etc/cron.hourly
/etc/cron.hourly/.placeholder
/etc/cron.monthly
/etc/cron.monthly/.placeholder
/etc/cron.weekly
/etc/cron.weekly/.placeholder
/etc/crontab
/etc/default
/etc/default/cron
/etc/init.d
/etc/init.d/cron
/etc/pam.d
/etc/pam.d/cron
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/cron.service
/usr
/usr/bin
/usr/bin/crontab
/usr/sbin
/usr/sbin/cron
/usr/share
/usr/share/bug
/usr/share/bug/cron
/usr/share/bug/cron/control
/usr/share/bug/cron/script
/usr/share/doc
/usr/share/doc/cron
/usr/share/doc/cron/FEATURES
/usr/share/doc/cron/NEWS.Debian.gz
/usr/share/doc/cron/README
/usr/share/doc/cron/README.Debian
/usr/share/doc/cron/README.anacron
/usr/share/doc/cron/THANKS
/usr/share/doc/cron/TODO.Debian
/usr/share/doc/cron/changelog.Debian.gz
/usr/share/doc/cron/changelog.gz
/usr/share/doc/cron/copyright
/usr/share/doc/cron/examples
/usr/share/doc/cron/examples/cron-stats.pl
/usr/share/doc/cron/examples/cron-tasks-review.sh
/usr/share/doc/cron/examples/crontab2english.pl
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/crontab.1.gz
/usr/share/man/man5
/usr/share/man/man5/crontab.5.gz
/usr/share/man/man8
/usr/share/man/man8/cron.8.gz
/var
/var/spool
/var/spool/cron
/var/spool/cron/crontabs
