# rootfs/var/lib/dpkg/info/dbus.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/session.d
/etc/dbus-1/system.d
/etc/default
/etc/default/dbus
/etc/init.d
/etc/init.d/dbus
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dbus.service
/lib/systemd/system/dbus.socket
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/sockets.target.wants
/usr
/usr/bin
/usr/bin/dbus-cleanup-sockets
/usr/bin/dbus-daemon
/usr/bin/dbus-monitor
/usr/bin/dbus-run-session
/usr/bin/dbus-send
/usr/bin/dbus-update-activation-environment
/usr/bin/dbus-uuidgen
/usr/lib
/usr/lib/dbus-1.0
/usr/lib/dbus-1.0/dbus-daemon-launch-helper
/usr/lib/sysusers.d
/usr/lib/sysusers.d/dbus.conf
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/dbus.conf
/usr/share
/usr/share/bug
/usr/share/bug/dbus
/usr/share/bug/dbus/control
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/session.conf
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system.conf
/usr/share/doc
/usr/share/doc/dbus
/usr/share/doc/dbus/AUTHORS.gz
/usr/share/doc/dbus/NEWS.gz
/usr/share/doc/dbus/README.Debian
/usr/share/doc/dbus/README.gz
/usr/share/doc/dbus/changelog.Debian.gz
/usr/share/doc/dbus/changelog.gz
/usr/share/doc/dbus/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/dbus
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dbus-cleanup-sockets.1.gz
/usr/share/man/man1/dbus-daemon.1.gz
/usr/share/man/man1/dbus-monitor.1.gz
/usr/share/man/man1/dbus-run-session.1.gz
/usr/share/man/man1/dbus-send.1.gz
/usr/share/man/man1/dbus-update-activation-environment.1.gz
/usr/share/man/man1/dbus-uuidgen.1.gz
/var
/var/lib
/var/lib/dbus
/lib/systemd/system/multi-user.target.wants/dbus.service
/lib/systemd/system/sockets.target.wants/dbus.socket
