# rootfs/var/lib/dpkg/info/haveged.conffiles
/etc/apparmor.d/usr.sbin.haveged
/etc/default/haveged
/etc/init.d/haveged
