# rootfs/var/lib/dpkg/info/ifupdown.conffiles
/etc/default/networking
/etc/init.d/networking
