# rootfs/var/lib/dpkg/info/ifupdown.list
/.
/etc
/etc/default
/etc/default/networking
/etc/init.d
/etc/init.d/networking
/etc/network
/etc/network/if-down.d
/etc/network/if-post-down.d
/etc/network/if-pre-up.d
/etc/network/if-up.d
/etc/network/interfaces.d
/lib
/lib/ifupdown
/lib/ifupdown/settle-dad.sh
/lib/ifupdown/wait-for-ll6.sh
/lib/ifupdown/wait-online.sh
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ifup@.service
/lib/systemd/system/ifupdown-pre.service
/lib/systemd/system/ifupdown-wait-online.service
/lib/systemd/system/networking.service
/lib/udev
/lib/udev/ifupdown-hotplug
/lib/udev/rules.d
/lib/udev/rules.d/80-ifupdown.rules
/sbin
/sbin/ifup
/usr
/usr/share
/usr/share/bug
/usr/share/bug/ifupdown
/usr/share/doc
/usr/share/doc/ifupdown
/usr/share/doc/ifupdown/NEWS.Debian.gz
/usr/share/doc/ifupdown/changelog.gz
/usr/share/doc/ifupdown/contrib
/usr/share/doc/ifupdown/contrib/ensureifup
/usr/share/doc/ifupdown/contrib/ifstate
/usr/share/doc/ifupdown/contrib/ifstate-check
/usr/share/doc/ifupdown/copyright
/usr/share/doc/ifupdown/examples
/usr/share/doc/ifupdown/examples/bridge
/usr/share/doc/ifupdown/examples/check-mac-address.sh
/usr/share/doc/ifupdown/examples/generate-interfaces.pl
/usr/share/doc/ifupdown/examples/get-mac-address.sh
/usr/share/doc/ifupdown/examples/network-interfaces
/usr/share/doc/ifupdown/examples/pattern-matching
/usr/share/doc/ifupdown/examples/pcmcia-compat.sh
/usr/share/doc/ifupdown/examples/ping-places.sh
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ifupdown
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/interfaces.5.gz
/usr/share/man/man8
/usr/share/man/man8/ifup.8.gz
/sbin/ifdown
/sbin/ifquery
/usr/share/man/man8/ifdown.8.gz
/usr/share/man/man8/ifquery.8.gz
