# rootfs/var/lib/dpkg/info/logrotate.conffiles
/etc/cron.daily/logrotate
/etc/logrotate.conf
/etc/logrotate.d/btmp
/etc/logrotate.d/wtmp
