# rootfs/var/lib/dpkg/info/openssh-server.list
/.
/etc
/etc/default
/etc/default/ssh
/etc/init.d
/etc/init.d/ssh
/etc/pam.d
/etc/pam.d/sshd
/etc/ssh
/etc/ssh/moduli
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/openssh-server
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rescue-ssh.target
/lib/systemd/system/ssh.service
/lib/systemd/system/ssh.socket
/lib/systemd/system/ssh@.service
/usr
/usr/lib
/usr/lib/openssh
/usr/lib/openssh/ssh-session-cleanup
/usr/sbin
/usr/sbin/sshd
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/openssh-server.py
/usr/share/doc
/usr/share/doc/openssh-client
/usr/share/doc/openssh-client/examples
/usr/share/doc/openssh-client/examples/ssh-session-cleanup.service
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/moduli.5.gz
/usr/share/man/man5/sshd_config.5.gz
/usr/share/man/man8
/usr/share/man/man8/sshd.8.gz
/usr/share/openssh
/usr/share/openssh/sshd_config
/usr/share/openssh/sshd_config.md5sum
/usr/share/doc/openssh-server
/usr/share/man/man5/authorized_keys.5.gz
