# rootfs/var/lib/dpkg/info/procps.conffiles
/etc/init.d/procps
/etc/sysctl.conf
/etc/sysctl.d/README.sysctl
/etc/sysctl.d/protect-links.conf
